# templates/etc/profile.d/go.sh.j2
export GOROOT={{ go_root }}
export GOPATH={{ go_path }}
export PATH=$PATH:${GOROOT}/bin:${GOPATH}/bin
